# rpmbuild/BUILD/sdobpkg-1.4/sdob-web.service
[Unit]
Description= python web hosting
PartOf=sdob.target

[Service]
Type=simple
#ExecStart=/usr/local/bin/python3.7 -m http.server 80 --directory /opt/sdob/
ExecStart=/bin/sh /opt/sdob/webserver.sh
Restart=on-failure

[Install]
WantedBy=sdob.target
