# rpmbuild/BUILD/sdobpkg-1.5/sdob-server.service
[Unit]
Description= sdob server
PartOf=sdob.target

[Service]
Type=simple
ExecStart=/usr/local/bin/python3.7 /opt/sdob/server.py
Restart=on-failure

[Install]
WantedBy=sdob.target
